# hamgr/container/etc/supervisord.conf
; A good starting point with documentation for this file can be generated
; from a supervisord installation by running echo_supervisord_conf

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:hamgr]
command=/etc/pf9-hamgr.sh
autostart=true
autorestart=true
startretries=99
startsecs=1
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
user=root
killasgroup=true
stopasgroup=true


[program:confd]
command=confd -backend consul -node %(ENV_CONFIG_HOST_AND_PORT)s -scheme %(ENV_CONFIG_SCHEME)s -auth-token %(ENV_CONSUL_HTTP_TOKEN)s -watch -prefix /customers/%(ENV_CUSTOMER_ID)s/regions/%(ENV_REGION_ID)s
redirect_stderr=true
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
autorestart=true
startretries=99
startsecs=1
