# hamgr/container/etc/supervisord.conf
; A good starting point with documentation for this file can be generated
; from a supervisord installation by running echo_supervisord_conf

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=debug
pidfile=/var/run/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:hamgr]
command=/etc/pf9-hamgr.sh
autostart=true
autorestart=true
startretries=99
startsecs=1
redirect_stderr=false
stderr_logfile=/var/log/pf9/hamgr/hamgr-err.log
stdout_logfile=/var/log/pf9/hamgr/hamgr.log
user=root
killasgroup=true
stopasgroup=true
